# dns/supervisor.conf
[program:dnschef]
command =   COPILOT_PLUGINS_DIRECTORYplugins/dns/start
directory = COPILOT_PLUGINS_DIRECTORYplugins/dns
user = root
autostart=true
autorestart=true
