# dns/supervisor.conf
[program:dnschef]
command =   %(ENV_COPILOT_PLUGINS_DIRECTORY)splugins/dns/start
directory = %(ENV_COPILOT_PLUGINS_DIRECTORY)splugins/dns
user = root
autostart=true
autorestart=true
