# suricata/supervisor.conf
[program:suricata]
command =   /home/www/copilot/copilot/plugins/suricata/start
directory = /home/www/copilot/copilot/plugins/suricata
user = root
autostart=true
autorestart=true
