# suricata/supervisor.conf
[program:suricata]
command =   COPILOT_PLUGINS_DIRECTORYplugins/suricata/start
directory = COPILOT_PLUGINS_DIRECTORYplugins/suricata
user = root
autostart=true
autorestart=true
