# create_ap/supervisor.conf
[program:create_ap]
command =  COPILOT_PLUGINS_DIRECTORYplugins/create_ap/start
directory = COPILOT_PLUGINS_DIRECTORYplugins/create_ap
user = root
autostart=true
autorestart=true
