# create_ap/supervisor.conf
[program:create_ap]
command =  %(ENV_COPILOT_PLUGINS_DIRECTORY)plugins/create_ap/start
directory = %(ENV_COPILOT_PLUGINS_DIRECTORY)plugins/create_ap
user = root
autostart=true
autorestart=true
